# systemd/remote-cli.service
[Unit]
Description=Secure Remote Chat CLI (Telegram)
Documentation=https://github.com/akamaor/remote-cli
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chatcli
Group=chatcli
WorkingDirectory=/opt/remote-cli

# Secrets loaded from root-owned file — token never appears in process table
EnvironmentFile=/opt/remote-cli/.env

ExecStart=/opt/remote-cli/venv/bin/python -m src.main

Restart=on-failure
RestartSec=15
StartLimitIntervalSec=60
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=remote-cli

# -----------------------------------------------------------------------
# Security note:
# This bot is a personal admin tool protected by WireGuard (no open inbound
# ports) + Telegram user-ID allowlist.  The systemd sandbox is intentionally
# minimal — heavy restrictions like NoNewPrivileges or ProtectSystem would
# block sudo, /home access, and arbitrary command execution, which are the
# whole point of the bot.
# -----------------------------------------------------------------------

# Kernel tunables and modules are off-limits for the Python process itself
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Restrict network families to what's needed (Telegram HTTPS + DNS + Unix)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# No real-time scheduling needed
RestrictRealtime=yes

# Prevent the Python process from changing CPU execution personality
LockPersonality=yes

# Clean up IPC objects when the service stops
RemoveIPC=yes

[Install]
WantedBy=multi-user.target
